# meshtastic-es-map-api.service
[Unit]
Description=Meshtastic-es-map API — sirve JSON desde SQLite al mapa web
After=network.target meshtastic-es-map-collector.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/meshtastic-es-map
ExecStart=/opt/meshtastic-es-map/venv/bin/python /opt/meshtastic-es-map/api/api_server.py --host 127.0.0.1 --port 8090 --db /opt/meshtastic-es-map/data/meshtastic-es-map.db
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
